# .local/services/xcompmgr/run
#!/bin/sh

if pidof xinit ;
then
export DISPLAY=:0
export XAUTHORITY=/run/user/1000/Xauthority

exec picom --config ~/.config/picom/picom.conf -r 15 -f
fi
