# .local/services/start-pulseaudio-x11/run
#!/bin/sh

export DISPLAY=:0
export XAUTHORITY=$HOME/.Xauthority

exec start-pulseaudio-x11
